# local/systemd/nipa-vmksft-native.service.sample
[Unit]
Description=Long-lived local vmksft service (native)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=#USER#
WorkingDirectory=#NIPA#/local
EnvironmentFile=#NIPA#/local/.env
ExecStart=/usr/bin/python3 #NIPA#/local/bin/vmksft_service.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
